# configs/launchd/com.user.research-daily-summary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.research-daily-summary</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/johannkaspar/Documents/Obsidian/Research/.scripts/run_daily_summary.sh</string>
  </array>

  <!-- Weekdays at 18:30 local time. Each day needs its own entry. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key><integer>1</integer>
      <key>Hour</key><integer>18</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>2</integer>
      <key>Hour</key><integer>18</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>3</integer>
      <key>Hour</key><integer>18</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>4</integer>
      <key>Hour</key><integer>18</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>5</integer>
      <key>Hour</key><integer>18</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
  </array>

  <!-- If the Mac was asleep at 18:30, run as soon as it wakes. -->
  <key>StartCalendarIntervalRunAtBoot</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/johannkaspar/Library/Logs/research-daily-summary.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/johannkaspar/Library/Logs/research-daily-summary.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
